# crontab
15 1 * * * ttrss cd /srv/www && git pull
*/30 * * * * ttrss /usr/bin/php7.0 /srv/www/update.php --feeds
